# includedir/etc/xdg/autostart/nm-tray.desktop
[Desktop Entry]
Exec=/bin/nm-tray
Name=nm-tray
OnlyShowIn=LXQt;
Type=Application
Version=1.0
X-LXQt-Need-Tray=true
